# scripts/vision.obrien.imagepipeline.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>vision.obrien.imagepipeline</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/justin/Documents/map/obrien.vision/scripts/run_image_pipeline.sh</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>

  <key>WatchPaths</key>
  <array>
    <string>/Users/justin/Documents/map/obrien.vision/uploads/new</string>
  </array>

  <key>StandardOutPath</key>
  <string>/Users/justin/Library/Logs/image_pipeline.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/justin/Library/Logs/image_pipeline.err.log</string>

  <key>WorkingDirectory</key>
  <string>/Users/justin/Documents/map/obrien.vision</string>
</dict>
</plist>